# debian/libapache2-mod-tile.prerm
#!/bin/sh
set -e

if [ "$1" = remove ] ; then
    a2dissite tileserver_site || true
    a2ensite 000-default || true
    a2dismod -q -f tile || true
    invoke-rc.d apache2 restart || true
fi

# automatically remove munin plugins
if [ -d /etc/munin/plugins ]
then
    rm -f /etc/munin/plugins/mod_tile* 
fi

#DEBHELPER#
