# debian/renderd.postinst
#!/bin/sh

mkdir -p /var/lib/mod_tile
chown www-data:www-data /var/lib/mod_tile

update-rc.d renderd defaults

